# user/vpn/vpn-notifier.timer
[Unit]
Description=Run VPN notifier periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
